# server/90-dac.rules
ACTION=="add", ATTR{idVendor}=="2a39", ATTR{idProduct}=="3fd3", SYMLINK+="dac", RUN+="/usr/local/bin/lsp-dac_added.sh"
